# scripts/bitm-ng.service
[Unit]
Description=BitM-NG Server
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/opt/bitm_ng
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

